# root/etc/systemd/system/ext_usb_reset.service
[Unit]
Description=gateway ext usb power cycle on start

[Service]
Type=oneshot
ExecStart=/data/p1/usb-pwr-cycle.py

[Install]
WantedBy=basic.target
